# deploy/uwsgi_upstart.conf
description "uwsgi loan_requests"
start on runlevel [2345]
stop on runlevel [06]

expect daemon

exec /srv/loan_requests/env/bin/uwsgi --ini /srv/loan_requests/project/deploy/uwsgi_site.ini

